# docker/image/config/common/ramips-mt7621/phicomm_k2p/etc/uci-defaults/99_custom
#!/bin/sh

set -e

uci -q batch <<-EOI >/dev/null
set system.@system[0].hostname='k2p'
commit system
EOI

exit 0
