# docker/image/config/21.02/ath79-nand/netgear_wndr4300/etc/uci-defaults/99_custom
#!/bin/sh

set -e

uci -q batch <<-EOI >/dev/null
set system.@system[0].hostname='wndr4300'
commit system
EOI

exit 0
